# buildroot/ovly/etc/udev/rules.d/99-usb_update.rules
ACTION=="add", KERNEL=="sd?[1-9]", SUBSYSTEMS=="usb", DRIVERS=="usb-storage", RUN+="/usr/local/bin/usb_update.sh %k"